# infra/provider-media-relay/provider-media-cleanup.timer
[Unit]
Description=Schedule ShanHaiEdu provider media expiry cleanup

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
RandomizedDelaySec=10s
Persistent=true
Unit=provider-media-cleanup.service

[Install]
WantedBy=timers.target
